# misc/systemd/treeland.service.in
[Unit]
Description=Treeland global wayland compositor service, only running on DDM.
PartOf=graphical.target
StartLimitIntervalSec=30
StartLimitBurst=2

Requires=seatd.service
After=seatd.service

[Service]
User=dde
Environment=LIBSEAT_BACKEND=seatd
Environment=DSG_APP_ID=org.deepin.dde.treeland
Environment=ASAN_SYMBOLIZER_PATH=/usr/bin/llvm-symbolizer
# For Debug: if enable asan, you can't get it's log from journalctl, so ensure the log to a file
Environment=ASAN_OPTIONS=log_path=/tmp/treeland-asan:detect_leaks=0:abort_on_error=1:symbolize=1
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/treeland.sh --lockscreen
Restart=on-failure
RestartSec=1s
# Disable logs from treeland.sh to avoid duplicate entries in journalctl
# since exec replaces the process, logs would appear from both treeland.sh and treeland
StandardOutput=null
StandardError=null

NoNewPrivileges=true
OOMScoreAdjust=-300
Nice=-15

# Enable MemoryDenyWriteExecute will cause Treeland crash in certain
# condition (observed on arm64 virtual machine). The coredump indicate
# that it happened in libgallium.so, which is involved by
# qwoutput->commit_state() call inside Output::enable().
#
# TODO: check if it is an internal mistake that used W^X memory there.
#
# MemoryDenyWriteExecute=true

PrivateIPC=true
ProtectSystem=full
ProtectHome=true
ProtectClock=true
ProtectProc=invisible
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true
