# data/systemd/dde-session-pre.target.wants/treeland-xwayland.service.in
[Unit]
Description=TreeLand socket helper service
StartLimitBurst=3
CollectMode=inactive-or-failed

Requires=treeland-xwayland.socket
PartOf=treeland-xwayland.socket
After=treeland-xwayland.socket

Before=dde-display.service

Requisite=dde-session-pre.target
PartOf=dde-session-pre.target
Before=dde-session-pre.target

[Service]
ExecCondition=/bin/sh -c 'test "$DDE_CURRENT_COMPOSITOR" = "TreeLand" || exit 2'
Type=notify
Sockets=treeland-xwayland.socket
UnsetEnvironment=DISPLAY
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/treeland-sd --type xwayland
ExecStopPost=-/usr/bin/systemctl --user unset-environment DISPLAY
Slice=session.slice
Restart=on-failure
RestartSec=1s
